# cloud/packer/startup-scripts/vllm.service
[Unit]
Description=vLLM Server
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/workspace
ExecStart=/usr/bin/bash /workspace/run_vllm.sh

[Install]
WantedBy=multi-user.target
